# config/crontab
0 2 * * * /app/bin/console r2d2:cleanup:room-availability
